# systemd/certbot-renew.timer
[Unit]
Description=Monthly TLS Certificate Renewal Timer

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
